# provisioning/files/etc/supervisor/conf.d/lou.conf
[program:lou]
command=start.sh
directory=/home/pi/lou/
user=pi
autostart=true
autorestart=true
stderr_logfile=/var/log/lou/err.log
stdout_logfile=/var/log/lou/out.log
