# provisioning/files/etc/supervisor/conf.d/lou.conf
[program:lou]
command=/home/pi/lou/start.sh
directory=/home/pi/lou
autostart=true
autorestart=true
stderr_logfile=/var/log/lou/err.log
stdout_logfile=/var/log/lou/out.log
