# provisioning/files/etc/supervisor/conf.d/lou.conf
[program:lou]
command=/home/pi/lou/npm run start-headless
autostart=true
autorestart=true
stderr_logfile=/var/log/lou-err.log
stdout_logfile=/var/log/lou-out.log
